# rootfs/root/crontabs/root
0 1 * * * /usr/local/bin/geoipupdate.sh
